# hostsoftware/network-autoconfig/debian/postinst
#!/bin/sh

set -e

if [ -x "/etc/init.d/net-choose" ]; then
	update-rc.d net-choose defaults >/dev/null
fi
