# hostsoftware/network-autoconfig/debian/postrm
#!/bin/sh

set -e

if [ "$1" = "purge" ] ; then
	update-rc.d net-choose remove >/dev/null
fi
