# hostsoftware/hexanode/backend/etc/service/hexanode-backend/run
#!/bin/sh

setuidgid hexanode hexabus_webbackend -u 127.0.0.1 -i eth0
