# hostsoftware/hexanode/webfrontend/etc/service/hexanode/run
#!/bin/sh

exec setuidgid hexanode server=127.0.0.1 port=3000 npm start hexanode-webfrontend
