# hostsoftware/libhexabus/etc/service/hexalog/run
#!/bin/sh

STORE=~/hexalog/store.db
IFACE=eth0

exec setuidgid hexalog hexalog -I $IFACE -s $STORE 2>&1
