# hostsoftware/hexanode/backend/etc/cron.d/hexabus_msg_heartbeat
CONFIG=/etc/hexabus_msg_bridge.conf
PATH=/usr/bin

*/10 * * * * root (hexabus_msg_bridge -c $CONFIG -H; echo $?) >/var/run/hexabus_msg_heartbeat.state 2>&1
