# hostsoftware/hexanode/webfrontend/debian/postinst
#!/bin/sh

set -e

if ! id hexanode >/dev/null 2>&1
then
	useradd -r -m hexanode
	update-rc.d reverse-ssh defaults
fi

su -c 'npm install /usr/share/hexabus_webfrontend' hexanode
