# hostsoftware/hexabridge/etc/service/hexabridge/run
#!/bin/sh

exec setuidgid hexabus hexabridge 2>&1
